# templates/systemd/openclaw-backup.timer
[Unit]
Description=Run OpenClaw backup daily

[Timer]
OnCalendar=__ONCALENDAR__
Persistent=true
Unit=openclaw-backup.service

[Install]
WantedBy=timers.target
